# sidecar-mirror-container/swat-agent.ini
[supervisord]
nodaemon=true

[program:swat-agent]
command=bash -c "source /home/root/secure-dir/swat.env && /usr/local/bin/scheduler"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=true
autorestart=true
startsecs=0
stopwaitsecs=1200
